# logwatch/logwatch.cron
# ------------------------------------------------------------------------
# THIS FILE HAS BEEN AUTOMATICALLY GENERATED VIA SALT-BASED CONFIGURATION.
# ANY MANUAL CHANGES WILL BE OVERWRITTEN!
# ------------------------------------------------------------------------

kernel_ignore_messages="tag#0 Sense Key : Hardware Error"
df_options="-h -l -x tmpfs -x overlay"

4 * * * * root	/usr/sbin/logwatch --output stdout --service zz-disk_space --service kernel --range 'since 1 hour ago for that hour'
